# root/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bash
exec cloudflared tunnel run
